# debian/postinst
#!/bin/bash -e

echo "Enabling MALOS wakeword service"
systemctl enable matrixio-malos-wakeword
echo "Please restart your Raspberry Pi after installation"
