# debian/prerm
#!/bin/bash -e

echo "Removing MALOS wakeword service"
systemctl stop matrixio-malos-wakeword
systemctl disable matrixio-malos-wakeword
systemctl daemon-reload
systemctl reset-failed


